# debian/postinst
#!/bin/bash
set -e

VENV=/usr/lib/akprx/venv
SERVICE_USER=akprx
SERVICE_GROUP=akprx

case "$1" in
    configure)
        # Create system group
        if ! getent group "$SERVICE_GROUP" > /dev/null 2>&1; then
            addgroup --system "$SERVICE_GROUP"
        fi

        # Create system user
        if ! getent passwd "$SERVICE_USER" > /dev/null 2>&1; then
            adduser \
                --system \
                --ingroup "$SERVICE_GROUP" \
                --no-create-home \
                --shell /usr/sbin/nologin \
                --gecos "akprx daemon" \
                "$SERVICE_USER"
        fi

        # Create directories
        mkdir -p /var/lib/akprx
        mkdir -p /var/lib/akprx/secrets

        # Set ownership and permissions
        chown "$SERVICE_USER:$SERVICE_GROUP" /var/lib/akprx
        chmod 750 /var/lib/akprx

        chown "root:$SERVICE_GROUP" /var/lib/akprx/secrets
        chmod 750 /var/lib/akprx/secrets

        # Create secrets file if missing
        if [ ! -f /var/lib/akprx/secrets/secrets.env ]; then
            touch /var/lib/akprx/secrets/secrets.env
        fi
        chown "$SERVICE_USER:$SERVICE_GROUP" /var/lib/akprx/secrets/secrets.env
        chmod 600 /var/lib/akprx/secrets/secrets.env

        # Create adaptors file if missing
        if [ ! -f /var/lib/akprx/adaptors.json ]; then
            echo '{}' > /var/lib/akprx/adaptors.json
        fi
        chown "$SERVICE_USER:$SERVICE_GROUP" /var/lib/akprx/adaptors.json
        chmod 600 /var/lib/akprx/adaptors.json

        # Install Python dependencies into venv
        if [ ! -d "$VENV" ]; then
            python3 -m venv "$VENV"
        fi
        "$VENV/bin/pip" install --quiet --upgrade pip
        "$VENV/bin/pip" install --quiet fastapi uvicorn httpx pydantic

        # Enable and start service
        systemctl daemon-reload
        systemctl enable --now akprx.service || true

        echo ""
        echo "akprx installed."
        echo ""
        echo "  akprx status"
        echo "  akprx key add <KEY_NAME>"
        echo "  akprx adaptor add"
        echo ""
    ;;
esac

#DEBHELPER#
